# uhk-usb-service/DEBIAN/postinst
#!/bin/bash
set -e

# If systemctl is available, reload units and enable the service now.
if command -v systemctl >/dev/null 2>&1; then
	# Reload systemd units to pick up the newly installed unit file
	systemctl daemon-reload || true
	# Enable and start the service for the current system
	systemctl enable --now uhk-usb-rebind.service || true
fi

exit 0
